# packages/finit/.footprint
drwxr-xr-x	root/root	etc/
-rw-r--r--	root/root	etc/rc.conf
drwxr-xr-x	root/root	etc/rc.d/
drwxr-xr-x	root/root	sbin/
-rwxr-xr-x	root/root	sbin/finit
lrwxrwxrwx	root/root	sbin/init -> finit
-rwxr-xr-x	root/root	sbin/reboot
